# infra/templates/app.service.example
[Unit]
Description=Audio Control App
After=network-online.target mopidy.service

[Service]
User=audio
WorkingDirectory=/home/audio/app/current/.next/standalone
ExecStart=/usr/bin/node apps/web/server.js
Restart=always
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=MOPIDY_RPC_URL=http://audio-os.local:6680/mopidy/rpc

[Install]
WantedBy=multi-user.target
